# profile/.config/profile.d/env.sh
# Shared environment variables that should work in POSIX shells.
# Fill this file with shell-agnostic exports only.
