# profile/.config/profile.d/linux.sh
# Linux-specific shared shell configuration belongs here.
# Keep this file POSIX-compatible because .profile sources it directly.
