# profile/.config/profile.d/local.sh
# Optional local machine overrides belong here.
# Source private files conditionally and avoid committing secrets directly.
