# profile/.config/profile.d/path.sh
# Shared PATH setup belongs here.
# Keep logic portable so it can be reused outside zsh.
